# crates/alex/config/systemd/alex.service
[Unit]
Description=Alex LLM proxy daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/alex daemon
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info,alex=debug
# systemd user services get a bare PATH (/usr/local/bin:/usr/bin), which hides
# harnesses installed under the user's home (claude, codex via mise, grok, …).
Environment=PATH=%h/.local/bin:%h/.local/share/mise/shims:%h/.grok/bin:%h/bin:/usr/local/bin:/usr/bin
NoNewPrivileges=true
ProtectSystem=full
# A fresh Linux user may not have used npm yet. The leading "-" keeps the
# optional cache path writable when present without making service startup
# fail when the directory does not exist.
ReadWritePaths=%h/.alex -%h/.npm
PrivateTmp=true
TimeoutStopSec=300

[Install]
WantedBy=default.target
