# crates/alex/config/systemd/alexandria.service
[Unit]
Description=Alex LLM proxy daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/alexandria daemon
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info,alexandria=debug
NoNewPrivileges=true
ProtectSystem=full
# A fresh Linux user may not have used npm yet. The leading "-" keeps the
# optional cache path writable when present without making service startup
# fail when the directory does not exist.
ReadWritePaths=%h/.alex -%h/.alexandria -%h/.npm
PrivateTmp=true
TimeoutStopSec=300

[Install]
WantedBy=default.target
